# import_photos/systemd.service
[Unit]
Description=import-photos-mail service

[Service]
User=root
ExecStart=/usr/local/bin/import-photos-mail.sh

[Install]
WantedBy=default.target
