# backups/systemd.service
[Unit]
Description=backups service

[Service]
User=root
ExecStart=/usr/local/bin/backup.sh

[Install]
WantedBy=default.target
