# access_service/access_service.service
[Unit]
Description=Allow or remove access

[Service]
ExecStart=/usr/bin/access_service.sh

[Install]
WantedBy=multi-user.target