# squid/systemd.service
[Unit]
Description=squidguard getlist service

[Service]
ExecStart=/usr/local/bin/getlists.sh

[Install]
WantedBy=default.target
